# scripts/insightpipe-import.service
[Unit]
Description=InsightPipe Import Service
After=network.target

[Service]
Type=simple
User=proxdoc
Group=proxdoc
WorkingDirectory=/opt/dockerapps/insightpipe
ExecStart=/bin/bash /opt/dockerapps/insightpipe/scripts/start-insightpipe-import.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (add as needed)
# Environment="VARIABLE_NAME=value"

[Install]
WantedBy=multi-user.target
